# src/main/scripts/deb/control/prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/arduieensy-extreme-feedback" ]; then
    invoke-rc.d arduieensy-extreme-feedback stop
fi
